# .config/autostart/at-spi-dbus-bus.desktop
[Desktop Entry]
Type=Application
Name=AT-SPI D-Bus Bus
Exec=/usr/lib/at-spi-bus-launcher --launch-immediately
Terminal=false
StartupNotify=false
