# .config/autostart/dunst.desktop
[Desktop Entry]
Exec=dunst
Name=dunst
Comment=A customizable and lightweight notification-daemon
Type=Application
Terminal=false
StartupNotify=false
