# .config/autostart/geoclue-demo-agent.desktop
[Desktop Entry]
Type=Application
Name=Geoclue Demo agent
Exec=/usr/lib/geoclue-2.0/demos/agent
NoDisplay=true
Terminal=false
StartupNotify=false
