# .config/autostart/gnome-keyring-pkcs11.desktop
[Desktop Entry]
Type=Application
Exec=gnome-keyring-daemon --start --components=pkcs11
NoDisplay=true
Terminal=false
StartupNotify=false
