# .config/autostart/gnome-keyring-ssh.desktop
[Desktop Entry]
Type=Application
Comment=GNOME Keyring: SSH Agent
Exec=gnome-keyring-daemon --start --components=ssh
NoDisplay=true
Terminal=false
StartupNotify=false
