# .config/autostart/picom.desktop
[Desktop Entry]
Version=1.0
Type=Application
NoDisplay=false
Name=picom
Categories=Utility;
Exec=/usr/bin/picom
StartupNotify=false
Terminal=false
