# .config/autostart/volumeicon.desktop
[Desktop Entry]
Exec=volumeicon
Name=volumeicon
Type=Application
Terminal=false
